# roles/common/templates/openshift.service
[Unit]
Description=OpenShift
After=docker.service
Requires=docker.service
Documentation=https://github.com/openshift/origin

[Service]
Type=simple
EnvironmentFile=-{{home_path}}openshift/origin
ExecStart={{origin_path}}/_output/local/bin/linux/amd64/openshift start $OPTIONS
LimitNOFILE=4096
WorkingDirectory={{home_path}}/

[Install]
WantedBy=multi-user.target
